# deploy/systemd/uo_regulations_translate_worker.service.template
[Unit]
Description=PDF OCR Translate Worker
After=network.target

[Service]
Type=simple
User={{APP_USER}}
Group={{APP_USER}}
WorkingDirectory={{APP_ROOT}}
Environment=PYTHONUNBUFFERED=1
EnvironmentFile={{ENV_FILE}}
ExecStart={{APP_ROOT}}/.venv/bin/python worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
