# deploy/systemd/uo_regulations_translate_log_cleanup.timer.template
[Unit]
Description=Run PDF OCR Translate Log Cleanup on schedule

[Timer]
OnCalendar={{CLEANUP_ON_CALENDAR}}
Persistent=true

[Install]
WantedBy=timers.target
