# deploy/systemd/uo_regulations_translate_template_backup.timer.template
[Unit]
Description=Run PDF OCR Translate Template Backup on schedule

[Timer]
OnCalendar={{TEMPLATE_BACKUP_ON_CALENDAR}}
Persistent=true

[Install]
WantedBy=timers.target
